# systemd/.config/systemd/user/dunst.service.d/override.conf
[Service]
Environment=DISPLAY=:0
After=xorg-server@0.service
