# systemd/.config/systemd/user/dunst.service.d/override.conf
[Unit]
Description=Dunst notification daemon
Documentation=man:dunst(1)
WantedBy=xorg.target

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
Environment=DISPLAY=:0
