# systemd/.config/systemd/user/bspwm.service
[Unit]
Descriptio=Binary Splitting Window Manager

[Service]
Type=simple
ExecStart=/usr/bin/bspwm

[Install]
WantedBy=xinit.target

