# systemd/.config/systemd/user/compton.service
[Unit]
Descriptio=Compositor

[Service]
Type=simple
ExecStart=/usr/bin/compton

[Install]
WantedBy=xinit.target

