# systemd/.config/systemd/user/compton.service
[Unit]
Description=Compositor

[Service]
Type=simple
ExecStart=/usr/bin/picom
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=xinit.target

