# systemd/.config/systemd/user/compton.service
[Unit]
Description=Compositor

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/picom
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=xorg.target

