# systemd/.config/systemd/user/kdeconnectd.service
[Unit]
Description=Connect your computer to your smartphone or tablet. Ripped from https://github.com/KDE/plasma-systemd-integration/
BindsTo=xorg.target

[Service]
Environment=DISPLAY=:0
Type=dbus
ExecStart=/usr/lib/kdeconnectd
BusName=org.kde.kdeconnect

[Install]
WantedBy=xorg.target

