# systemd/.config/systemd/user/polybar@.service
[Unit]
Description=Bar
Wants=bspwm.service
After=bspwm.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/polybar -c %h/.config/polybar/config.ini %i

[Install]
WantedBy=xorg.target

