# systemd/.config/systemd/user/polybar@.service
[Unit]
Descriptio=Compositor

[Service]
Type=simple
ExecStart=/usr/bin/polybar -c %h/.config/polybar/config.ini %i

[Install]
WantedBy=xinit.target

