# systemd/.config/systemd/user/pywal.service
[Install]
Description=PyWal Reload Oneshot

[Service]
Type=oneshot
ExecStart=/usr/bin/wal -R

[Install]
WantedBy=xinit.target

