# systemd/.config/systemd/user/pywal.service
[Unit]
Description=PyWal Reload Oneshot

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/wal -R

[Install]
WantedBy=xorg.target

