# systemd/.config/systemd/user/pywal.service
[Unit]
Description=PyWal Reload Oneshot
BindsTo=xorg.target

[Service]
Environment=DISPLAY=:0
Type=oneshot
ExecStart=/usr/bin/wal -R

[Install]
WantedBy=xorg.target

