# systemd/.config/systemd/user/pywal.service
[Unit]
Description=PyWal Reload Oneshot
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wal -R

[Install]
WantedBy=graphical-session.target

