# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
PartOf=bspwm.service
After=bspwm.service
BindsTo=xorg.target

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/sxhkd
Environment="SXHKD_SHELL=/usr/bin/zsh"

[Install]
WantedBy=xorg.target

