# systemd/.config/systemd/user/sxhkd.service
[Unit]
Descriptio=Simple X Hotkey Daemon

[Service]
Type=simple
ExecStart=/usr/bin/sxhkd

[Install]
WantedBy=xinit.target

