# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
PartOf=bspwm.service

[Service]
Type=simple
ExecStart=/usr/bin/sxhkd
Environment="SXHKD_SHELL=/usr/bin/zsh"

[Install]
WantedBy=xinit.target

