# systemd/.config/systemd/user/xorg.target
[Unit]
Description=xorg systemd --user target
Requires=xorg-server@0.service

[Install]
WantedBy=basic.target
